# pagenotif.cron
#!/bin/sh

# place in /etc/cron.hourly for hourly execution
# place in /etc/cron.daily for daily execution
# place in /etc/cron.weekly for weekly execution

/usr/bin/python /usr/local/bin/pagenotif /etc/pagenotif.conf
